# sandbox/supervisord.conf
; 通信http服务配置
[unix_http_server]
file=/tmp/supervisor.sock    ; socket文件路径
chmod=0770                   ; socket文件权限，770表示文件所有者和同组用户有读/写/执行权限，其他用户无任何权限
chown=root:root              ; API进程以root运行；避免沙箱在drop-all能力下执行chown

; supervisord日志配置信息
[supervisord]
user=root                    ; 显式声明，消除root运行告警
logfile=/dev/stdout          ; 主日志路径
logfile_maxbytes=0           ; 不限制日志大小(用完即弃无需考虑)
loglevel=info                ; 日志等级
pidfile=/tmp/supervisord.pid ; pid文件路径
nodaemon=true                ; 设为true表示在前台运行(Docker前台)，方便Docker捕捉日志
minfds=1024                  ; 打开文件描述符的最小值
minprocs=200                 ; 最小进程数
autoshutdown=true            ; 所有服务停止后自动退出

; 启用supervisor的远程过程调用接口
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; supervisorctl连接服务配置
[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

; FastAPI进程配置
[program:app]
command=uvicorn app.main:app --host 0.0.0.0 --port 8080 %(ENV_UVI_ARGS)s
directory=/sandbox
user=root
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME=/home/ubuntu
priority=60

; Chrome浏览器配置
[program:chrome]
command=chromium \
    --display=:1 \
    --window-size=1280,1080 \
    --start-maximized \
    --no-sandbox \
    --disable-dev-shm-usage \
    --disable-setuid-sandbox \
    --disable-accelerated-2d-canvas \
    --disable-gpu \
    --disable-features=WelcomeExperience,SigninPromo \
    --no-first-run \
    --no-default-browser-check \
    --disable-infobars \
    --test-type \
    --disable-popup-blocking \
    --disable-gpu-sandbox \
    --no-xshm \
    --new-window=false \
    --disable-notifications \
    --disable-extensions \
    --disable-component-extensions-with-background-pages \
    --disable-prompt-on-repost \
    --disable-dialogs \
    --disable-modal-dialogs \
    --disable-web-security \
    --disable-site-isolation-trials \
    --remote-debugging-address=0.0.0.0 \
    --remote-debugging-port=8222 %(ENV_CHROME_ARGS)s
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=DISPLAY=:1
priority=20
startretries=3
startsecs=5

; socat进程配置
[program:socat]
command=socat TCP-LISTEN:9222,bind=0.0.0.0,fork,reuseaddr TCP:127.0.0.1:8222
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=30
startsecs=2

; Xvfb虚拟显示器配置
[program:xvfb]
command=bash -c "rm -f /tmp/.X1-lock && Xvfb :1 -screen 0 1280x1080x24"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10
startsecs=2

; x11vnc服务配置
[program:x11vnc]
command=x11vnc -display :1 -nopw -shared -listen 0.0.0.0 -xkb -forever -rfbport 5900
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=40
startsecs=3

; Websockify配置/将VNC转换为Websocket
[program:websockify]
command=websockify 0.0.0.0:5901 localhost:5900
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=50
startsecs=3

; 使用分组来统一启动or关闭进程
[group:services]
programs=xvfb,chrome,socat,x11vnc,websockify,app
